# supervisord.conf
; supervisord.conf — Process manager for the Charon container.
;
; Runs three processes inside a single Docker container:
;   clawrouter      — x402-native LLM router (internal port 8402)
;   relay           — Node.js TypeScript relay server (public port 3000)
;   openclaw-bridge — Python inference bridge (internal port 4001)
;
; Start order: clawrouter → relay → openclaw-bridge
; openclaw-bridge connects to both relay (:3000) and clawrouter (:8402).
;
; supervisord itself runs in the foreground (nodaemon=true) so Docker tracks it
; as PID 1. All child processes are restarted automatically if they crash.
;
; Environment variables are inherited from the container environment — Docker
; passes them in via `--env` or `--env-file`. Do NOT hard-code secrets here.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; ── ClawRouter — x402 LLM proxy ───────────────────────────────────────────────
; Starts first. Derives BLOCKRUN_WALLET_KEY from WDK_SEED_PHRASE via the wrapper.
; Port 8402 is internal — never exposed outside the container.
[program:clawrouter]
command=node /app/clawrouter-start.js
directory=/app
environment=BLOCKRUN_PROXY_PORT="8402"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Node.js relay ─────────────────────────────────────────────────────────────
[program:relay]
command=node /app/dist/index.js
directory=/app
; PORT and all WDK_* vars come from the container environment.
; RELAY_INTERNAL_URL tells the bridge where to find the internal API.
environment=RELAY_INTERNAL_URL="http://localhost:3000/internal"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Python OpenClaw bridge ────────────────────────────────────────────────────
[program:openclaw-bridge]
command=python3 /app/openclaw-bridge/server.py
directory=/app
environment=RELAY_INTERNAL_URL="http://localhost:3000/internal",BRIDGE_PORT="4001",CLAW_ROUTER_URL="http://localhost:8402/v1",MEMORY_PATH="/app/data/MEMORY.md"
autostart=true
autorestart=true
startretries=5
; Wait 5 seconds so both relay and clawrouter are up first.
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
